# application/denoapp.service
[Unit]
Description=Deno Todo Web Server
After=network.target

[Service]
Type=simple
User=denoapp
Group=denoapp
WorkingDirectory=/opt/webserver
ExecStart=/opt/webserver/todo-server
Restart=on-failure

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true

MemoryMax=512M
CPUWeight=50
IOWeight=10
TasksMax=50

[Install]
WantedBy=multi-user.target
