# auditd/hardening.rules
# Monitor /etc/passwd and /etc/shadow for changes
-w /etc/passwd -p wa -k passwd_changes
-w /etc/shadow -p wa -k shadow_changes

# Monitor user/group changes
-w /etc/group -p wa -k group_changes
-w /etc/gshadow -p wa -k gshadow_changes

# Monitor sudo commands
-w /var/log/sudo.log -p rwxa -k sudo_log

# Monitor nginx configuration changes
-w /etc/nginx/ -p wa -k nginx_conf

# Monitor execution of binaries in /bin and /usr/bin
-w /bin/ -p x -k bin_exec
-w /usr/bin/ -p x -k usrbin_exec

# 1. Critical Configuration Files
# SSH Configuration
-w /etc/ssh/sshd_config -p wa -k ssh_config

# PAM (Pluggable Authentication Modules) Configuration
-w /etc/pam.d/ -p wa -k pam_config

# Crontab and Scheduled Jobs
-w /etc/crontab -p wa -k crontab_changes
-w /etc/cron.d/ -p wa -k cron_d_changes
-w /var/spool/cron/crontabs/ -k user_crontabs

# Kernel Modules
-w /etc/modprobe.d/ -p wa -k modprobe_config
-a always,exit -F arch=b64 -S finit_module,init_module,delete_module -k kernel_module_changes
-a always,exit -F arch=b32 -S finit_module,init_module,delete_module -k kernel_module_changes_32bit

# Network Configuration
-w /etc/network/interfaces -p wa -k network_config
-w /etc/netplan/ -p wa -k netplan_config # For newer Ubuntu versions using netplan
-w /etc/iptables/rules.v4 -p wa -k iptables_v4
-w /etc/iptables/rules.v6 -p wa -k iptables_v6
-w /etc/ufw/ -p wa -k ufw_config

# 2. System Call Monitoring for Specific Events
# Mount Operations
-a always,exit -F arch=b64 -S mount,umount2 -F auid!=-1 -k mount_operations
-a always,exit -F arch=b32 -S mount,umount,umount2 -F auid!=-1 -k mount_operations_32bit

# System Time Changes
-a always,exit -F arch=b64 -S adjtimex,settimeofday,clock_settime -k time_change
-a always,exit -F arch=b32 -S adjtimex,settimeofday,clock_settime -k time_change_32bit
-w /etc/localtime -p wa -k localtime_change

# Process ID (PID) Changes (fork/clone) - Potentially noisy, review if necessary
-a always,exit -F arch=b64 -S fork,vfork,clone -k process_creation
-a always,exit -F arch=b32 -S fork,vfork,clone -k process_creation_32bit

# Privilege Escalation Attempts (execve with altered permissions)
-a always,exit -F arch=b64 -S execve -F euid=0 -k exec_root
-a always,exit -F arch=b32 -S execve -F euid=0 -k exec_root_32bit

# File Deletion
-a always,exit -F arch=b64 -S unlink,unlinkat,rename,renameat -k file_delete_rename
-a always,exit -F arch=b32 -S unlink,unlinkat,rename,renameat -k file_delete_rename_32bit

# 3. User Activity and Authentication
# Login/Logout Events (additional to syslog)
-w /var/log/faillog -p wa -k login_failure
-w /var/log/lastlog -p wa -k lastlog_update
-w /var/log/tallylog -p wa -k login_tally

# User and Group Modification Tools
-w /usr/sbin/useradd -p x -k user_add
-w /usr/sbin/usermod -p x -k user_mod
-w /usr/sbin/userdel -p x -k user_del
-w /usr/sbin/groupadd -p x -k group_add
-w /usr/sbin/groupmod -p x -k group_mod
-w /usr/sbin/groupdel -p x -k group_del
-w /usr/bin/passwd -p x -k passwd_cmd

# 4. Important Directories
# /var/log/ - Monitor for log tampering
-w /var/log/audit/ -p wa -k audit_log_changes
-w /var/log/apt/ -p wa -k apt_log_changes
-w /var/log/auth.log -p wa -k auth_log_changes

# /boot/ - Changes to boot loader or kernel images
-w /boot/ -p wa -k boot_changes

# /opt/ and /srv/ - Application/service directories
-w /opt/ -p wa -k opt_dir_changes
-w /srv/ -p wa -k srv_dir_changes